# etc/cron.d/ViOfflineInventory
3 *  * * * root /usr/bin/perl /root/getInventory.pl