# etc/cron.d/ViOfflineInventory
3 *  * * * root /usr/bin/pwsh -NonInteractive -NoProfile -f /opt/sexigraf/ViOfflineInventory.ps1 -credstore /var/www/.vmware/credstore/vipscredentials.xml >/dev/null 2>&1